# deploy/blog-aggregator.service
[Unit]
Description=Blog Aggregator API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/blog-aggregator
Environment="ADMIN_PASSWORD=your-secure-password"
Environment="OPENAI_API_KEY=your-openai-key"
Environment="SECRET_KEY=your-random-secret-key"
ExecStart=/opt/blog-aggregator/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
