# deploy/blog-aggregator.service
[Unit]
Description=Blog Aggregator API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/blog-aggregator
Environment="ADMIN_PASSWORD=change-me"
Environment="SECRET_KEY=change-me-random-string"
ExecStart=/opt/blog-aggregator/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
